# fix-it-fred.service
[Unit]
Description=Fix It Fred AI Assistant
Documentation=https://chatterfix.com
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/yoyofred_gringosgambit_com/chatterfix-docker/app/fix_it_fred_ai_service.py
User=yoyofred_gringosgambit_com
Group=yoyofred_gringosgambit_com
WorkingDirectory=/home/yoyofred_gringosgambit_com/chatterfix-docker/app
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=15
Environment=PORT=9000
Environment=OLLAMA_HOST=localhost:11434
Environment=DEFAULT_MODEL=llama3.2:1b

[Install]
WantedBy=multi-user.target
