# fix-it-fred-devops.service
[Unit]
Description=Fix It Fred DevOps - Autonomous VM Management
Documentation=https://github.com/TheGringo-ai/Chatterfix
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=yoyofred_gringosgambit_com
Group=yoyofred_gringosgambit_com
WorkingDirectory=/home/yoyofred_gringosgambit_com/chatterfix-docker
Environment=PYTHONPATH=/home/yoyofred_gringosgambit_com/chatterfix-docker
Environment=PYTHONUNBUFFERED=1
Environment=PORT=9003
Environment=FRED_MODE=production
Environment=FRED_LOG_LEVEL=INFO
Environment=FRED_HEALING_INTERVAL=60
Environment=FRED_DEPLOYMENT_CHECK_INTERVAL=300
Environment=GITHUB_REPO=TheGringo-ai/Chatterfix
Environment=GITHUB_BRANCH=main-clean

# Main Fred DevOps process
ExecStart=/usr/bin/python3 /home/yoyofred_gringosgambit_com/chatterfix-docker/fix_it_fred_devops_daemon.py

# Health check and restart logic
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Process management
OOMScoreAdjust=-100
Nice=-5

# Security settings (relaxed for development access)
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fix-it-fred-devops

[Install]
WantedBy=multi-user.target