# deployment-package/chatterfix-integrated/chatterfix-integrated.service
[Unit]
Description=ChatterFix Integrated CMMS System
After=network.target

[Service]
Type=forking
User=yoyofred_gringosgambit_com
WorkingDirectory=/home/yoyofred_gringosgambit_com/chatterfix-integrated
ExecStart=/home/yoyofred_gringosgambit_com/chatterfix-integrated/start_all_services.sh
ExecStop=/home/yoyofred_gringosgambit_com/chatterfix-integrated/stop_all_services.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
