# core/cmms/chatterfix-daily-summary.service
[Unit]
Description=ChatterFix Daily Summary Generator
Documentation=https://github.com/chatterfix/cmms

[Service]
Type=oneshot
User=chatterfix
Group=chatterfix
WorkingDirectory=/opt/chatterfix/core/cmms
Environment=PYTHONPATH=/opt/chatterfix/core/cmms
Environment=AI_MEMORY_DIR=/opt/chatterfix/ai-memory
ExecStart=/usr/bin/python3 -c "
import asyncio
from datetime import date, timedelta
from aiops_agent import aiops_agent

async def generate_summary():
    yesterday = (date.today() - timedelta(days=1))
    summary = await aiops_agent._generate_daily_summary(yesterday)
    
    summary_file = aiops_agent.lessons_dir / f'daily_summary_{yesterday.strftime(\"%Y%m%d\")}.md'
    with open(summary_file, 'w') as f:
        f.write(summary)
    
    print(f'Daily summary generated: {summary_file}')

asyncio.run(generate_summary())
"
TimeoutStartSec=300

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chatterfix/ai-memory
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chatterfix-daily-summary