# infra/config/cmms/chatterfix-cmms.service
[Unit]
Description=ChatterFix CMMS API (atomic releases via /opt/chatterfix-cmms/current)
After=network-online.target
Wants=network-online.target

[Service]
# Run as a dedicated user (create it once: useradd --system --no-create-home --shell /usr/sbin/nologin chatterfix)
User=chatterfix
Group=chatterfix

# Point at the symlinked release
WorkingDirectory=/opt/chatterfix-cmms/current/core/cmms

# Load env without crashing if missing (we'll still have sane defaults)
EnvironmentFile=-/etc/default/chatterfix-cmms

# Start via uvicorn in the release's venv; bind to env port with fallback
ExecStart=/opt/chatterfix-cmms/current/.venv/bin/uvicorn app:app --host 0.0.0.0 --port ${CMMS_PORT:-8000} --timeout-keep-alive 30

# Gate readiness: only mark started if /ready goes green within ~30s
ExecStartPost=/bin/bash -lc 'for i in {1..15}; do curl -fsS http://127.0.0.1:${CMMS_PORT:-8000}/ready >/dev/null && exit 0; sleep 2; done; exit 1'

Restart=always
RestartSec=2
TimeoutStartSec=60
TimeoutStopSec=20

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
# App state lives here; systemd creates it owned by chatterfix:chatterfix
StateDirectory=chatterfix-cmms
# Explicitly allow writes to /var/lib/chatterfix-cmms (needed with ProtectSystem=full)
ReadWritePaths=/var/lib/chatterfix-cmms
LimitNOFILE=65535

# Logs to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chatterfix-cmms

[Install]
WantedBy=multi-user.target