# infra/config/cmms/chatterfix-health-check.service
[Unit]
Description=ChatterFix Health Check (One-shot)
Documentation=https://github.com/chatterfix/cmms

[Service]
Type=oneshot
User=chatterfix
Group=chatterfix
WorkingDirectory=/opt/chatterfix/core/cmms
Environment=PYTHONPATH=/opt/chatterfix/core/cmms
ExecStart=/opt/chatterfix/scripts/health-check.sh
TimeoutStartSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp/chatterfix-health /opt/chatterfix/ai-memory
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chatterfix-health-check